# services/prplmesh/files/etc/uci-defaults/prplmesh

# nothing to do if there is already a prplmesh config:
uci -q get prplmesh && exit 0

touch /etc/config/prplmesh

uci batch <<'EOF'
set prplmesh.config='prplmesh'
set prplmesh.config.management_mode='Multi-AP-Agent'
set prplmesh.config.certification_mode='0'

commit prplmesh
EOF
